# deploy/systemd/minecraft-exporter.service
# =============================================================================
# Minecraft Prometheus Metrics Exporter - Systemd Service File
# =============================================================================
# Exposes Minecraft server metrics on port 9225 for Prometheus scraping
#
# Installation:
#   Use: sudo ./deploy/install-service.sh
#
# This is a TEMPLATE file. Placeholders are replaced during installation:
#   - DEPLOY_USER / DEPLOY_GROUP -> actual user running the service
# =============================================================================

[Unit]
Description=Minecraft Prometheus Metrics Exporter
Documentation=file:///opt/minecraft/CLAUDE.md
After=network.target
After=minecraft.service
BindsTo=minecraft.service

# Start after minecraft service
PartOf=minecraft.service

[Service]
Type=simple

# User and group - replaced during installation by install-service.sh
User=DEPLOY_USER
Group=DEPLOY_GROUP

# Working directory
WorkingDirectory=/opt/minecraft

# Environment file
EnvironmentFile=-/opt/minecraft/deploy/config.env

# Main exporter process
ExecStart=/bin/bash /opt/minecraft/deploy/prometheus/minecraft-exporter.sh

# Restart configuration
Restart=always
RestartSec=10

# Resource limits (exporter is lightweight)
MemoryMax=50M
CPUQuota=10%

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/opt/minecraft
NoNewPrivileges=yes
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=minecraft-exporter

[Install]
WantedBy=multi-user.target
